# systemd/cert_renew.service
[Unit]
Description=Renew all certificates previously obtained with certbot.

[Service]
ExecStart=ut4ms cert renew
