# systemd/cert_renew.service
[Unit]
Description=Certificate renewal

[Service]
ExecStart=ut4ms cert renew
