# systemd/cert_renew.service
[Unit]
Description=Certificate renewal

[Service]
Type=exec
ExecStart=ut4ms cert renew
