# systemd/logs_clear.service
[Unit]
Description=Clear old outdated logs.

[Service]
Type=exec
ExecStart=journalctl --vacuum-time=30d
