# systemd/system_upgrade.service
[Unit]
Description=Upgrade all packages on the system and reboot.

[Service]
ExecStart=dnf -q -y upgrade && update-certbot && update-docker-compose && reboot
