# systemd/ut4ms_update.service
[Unit]
Description=Try to update and rebuild UT4MasterServer

[Service]
ExecStart=ut4ms server update && systemctl reload ut4ms
